# eh_remoteaccess/rootfs/etc/s6-overlay/s6-rc.d/cloudflared/run
#!/command/with-contenv bashio
# ==============================================================================
# Home Assistant Add-on: Cloudflared
# Runs the Cloudflare Tunnel for Home Assistant
# ==============================================================================
declare config_file="/tmp/config.json"
declare certificate="/data/cert.pem"
declare -a options


# Set common cloudflared tunnel options
options+=(--no-autoupdate)
options+=(--metrics="0.0.0.0:36500")


# Check for additional run parameters
if bashio::config.has_value 'run_parameters' ; then
    bashio::log.trace "bashio::config.has_value 'run_parameters'"
    for run_parameter in $(bashio::config 'run_parameters'); do
        bashio::log.trace "Adding run_parameter: ${run_parameter}"
        options+=("${run_parameter}")
    done
fi

# Get the email address and system ID from the add-on configuration
EH_EMAIL_ADDRESS=$(bashio::config 'eh_email_address')
EH_SYSTEM_ID=$(bashio::config 'eh_system_id')

# Construct the URL for the API call
API_URL="https://ehsysteminitialize.effortlesshome.co/getremoteaccessinfo/${EH_EMAIL_ADDRESS}/${EH_SYSTEM_ID}"

# Log the URL for debugging (optional)
bashio::log.trace "Calling EH API at: ${API_URL}"

# Make the API call and store the response
RESPONSE=$(curl -s -X GET "${API_URL}")

# Check if the response contains valid data
if bashio::var.has_value "${RESPONSE}"; then
    bashio::log.info "EH API call successful, response: ${RESPONSE}"

    # Parse the JSON response using jq to extract Token and URL
    TOKEN=$(echo "${RESPONSE}" | jq -r '.Token')
    REMOTE_URL=$(echo "${RESPONSE}" | jq -r '.URL')

    bashio::log.trace "token retrieved"
    options+=(run --token="$TOKEN")
fi

bashio::log.info "Connecting Cloudflare Tunnel..."
bashio::log.debug "cloudflared tunnel ${options[@]}"
exec cloudflared \
        tunnel "${options[@]}"
